# install/etc/cron.d/mynotes-stats
# MyNotes graphics update script
*/5 * * * *	 root /home/inexika/bin/mynotes-stats.sh
